# docker/scripts/sb-openclaw.service
[Unit]
Description=2nd-brain OpenClaw gateway (compose up -d)
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/projects/2nd-brain-docker
ExecStart=/usr/bin/docker compose -f compose.openclaw.yml up -d
ExecStop=/usr/bin/docker compose -f compose.openclaw.yml down
TimeoutStartSec=180

[Install]
WantedBy=default.target
