# docker/parser-drain/parser-drain.timer
[Unit]
Description=2nd-brain parser drain — periodic poll

[Timer]
OnBootSec=2min
OnUnitInactiveSec=10min      # 드레인 *종료* 10분 후 재발화 (긴 드레인과 겹침 없음; 도는 서비스는 또 안 띄움)
Persistent=true

[Install]
WantedBy=timers.target
