# docker/parser-drain/parser-drain.timer
[Unit]
Description=2nd-brain parser drain — periodic poll

[Timer]
# OnActiveSec: 타이머 활성(부팅·enable·daemon-reload)마다 시드 → 체인 재무장.
# OnBootSec 단독은 부팅 때만 시드라 daemon-reload 후 NEXT 가 비어 재무장 실패함.
# (systemd 는 인라인 주석 미지원 — 주석은 반드시 별도 줄.)
OnActiveSec=2min
# OnUnitInactiveSec: 드레인 *종료* 10분 후 재발화 (긴 드레인과 겹침 없음; 도는 서비스는 또 안 띄움).
OnUnitInactiveSec=10min

[Install]
WantedBy=timers.target
